# php_scripts.ini
[program:php_scripts]
command = php /www/wwwroot/muamaster-tag/adminmua/think CalculateUserStatsByIdCommand
priority=6
autorestart = true
stderr_logfile = /data/logs/php_scripts.err
stderr_logfile_maxbytes = 200MB
stderr_logfile_backups = 10000
stdout_logfile = /data/logs/php_scripts.log
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 100000
